# app-admin/cmastor/files/cmastor
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later

#NB: Config is in /etc/conf.d/bonding

depend() {
	need net snmpd hpasm
}

start() {
	ebegin "Starting CMASTOR"
	/opt/compaq/storage/etc/cmastor start
	eend $?
}

stop() {
	ebegin "Stopping CMASTOR"
	/opt/compaq/storage/etc/cmastor stop
	eend $?
}

